# templates/sub-ca-sign-certs.path
[Unit]
Description= Triggers the service that signs CSRs.

[Path]
PathModified=<%= @certs_requests %>
PathModified=<%= @certs_server %>
PathModified=<%= @certs_client %>

[Install]
WantedBy=multi-user.target

