# templates/sub-ca-sign-certs.service
[Unit]
Description= Starts the cert signing script

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/sub-ca-sign-certs.sh

